# rootfs/etc/inittab
::restart:/sbin/init
::shutdown:sync
::shutdown:umount -a
::ctrlaltdel:/sbin/reboot
::respawn:/etc/custominit